# src/fty-config.service.in
[Unit]
Description=fty-config service
After=network.target malamute.service
Requires=malamute.service network.target
PartOf=bios.target

[Service]
Type=simple
User=root
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-config --config @sysconfdir@/@PACKAGE@/fty-config.cfg
Restart=always

[Install]
WantedBy=bios.target
